# files/etc/cron.d/decent.im_letsencrypt
# m hh dom mon dow user       command
0   3   1   *   *  root       /sbin/decent.im_letsencrypt > /var/log/decent.im_letsencrypt.log 2> /var/log/decent.im_letsencrypt.err
