# files/etc/cron.d/decent.im_update_world
# m hh dom mon dow user       command
0   12  *   *  tue root       bash -l /sbin/decent.im_update_world > /var/log/decent.im_update_world.log 2> /var/log/decent.im_update_world.err

# Assuming that update may harm system operational status, let's perform update
# during working week & time so that admin is more probably able to notice and
# fix the issue.
