# files/etc/cron.d/jabber_backup
# m hh dom mon dow user       command
0   0   *   *   *  root       /usr/local/sbin/jabber_db_backup_daily.sh
