# files/etc/cron.daily/decent.im_health_check
#!/bin/sh
exec chronic decent.im_health_check
