# files/etc/profile.d/nullmailer-to.sh
# nullmailer-inject(1):
# Insert a To line containing a list of the recipients if the header does not contain either a To...
export NULLMAILER_FLAGS=t
