# files/etc/cron.weekly/decent.im_update_netfilter_rules
#!/bin/bash
exec decent.im_update_netfilter_rules
