# files/etc/cron.weekly/decent.im_cron_mail_test
#!/bin/sh
echo "Weekly test of mail from cron"
